# systemd/backup.service
[Unit]
Description=Restic backup service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/restic-backup.sh
